# supervisord/restart.ini
[program:php127-msgaudit-restart]
process_name=%(program_name)s_%(process_num)02d
directory=/www/wwwroot/php127-msgaudit/bin/
command=sh /www/wwwroot/php127-msgaudit/bin/restart.sh
autostart=true
autorestart=true
user=root
numprocs=1
redirect_stderr=true
stdout_logfile=/www/wwwroot/php127-msgaudit/supervisord/logs/restart.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=10